# overlay/ec/lib/svc/manifest/ec-autostart.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="ec-autostart">
	<service name="ec-autostart" type="service" version="1">
		<create_default_instance enabled="true" />
		<single_instance />
		<dependency name='network' grouping='require_any' restart_on='error' type='service'>
		    <service_fmri value='svc:/milestone/network'/>
		</dependency>
		<method_context>
			<method_credential user="root" group="root" />
		</method_context>
		<exec_method type="method" name="start" exec="/lib/svc/method/ec-autostart"
			timeout_seconds="60" />
		<exec_method type="method" name="stop" exec=":kill"
			timeout_seconds="60" />
		<property_group name="startd" type="framework">
			<propval name="duration" type="astring" value="transient" />
			<propval name="ignore_error" type="astring" value="core,signal" />
		</property_group>
		<property_group name="application" type="application">
		</property_group>
		<stability value="Evolving" />
		<template>
			<common_name>
				<loctext xml:lang="C">
					EveryCity startup script
				</loctext>
			</common_name>
		</template>
	</service>
</service_bundle>

